# rootfs/etc/services.d/vue/run
#!/usr/bin/with-contenv sh

cd /var/www;

FILE=package.json
if [ -f "$FILE" ]; then
    exec npm run serve
fi
